# scripts/online-judge-backup.service
[Unit]
Description=JudgeKit PostgreSQL backup
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/judgekit
ExecStart=/bin/bash -lc 'mkdir -p /home/ubuntu/backups/judgekit && DB_DIALECT=postgresql /bin/bash /home/ubuntu/judgekit/scripts/backup-db.sh "/home/ubuntu/backups/judgekit/judge-$(date +%%Y%%m%%d-%%H%%M%%S).sql.gz"'
