# scripts/online-judge-worker-rs.service
[Unit]
Description=Online Judge Worker (Rust)
After=network.target docker.service online-judge.service
Requires=online-judge.service
OnFailure=notify-failure@%n.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/judgekit
EnvironmentFile=/home/ubuntu/judgekit/.env
Environment=DEAD_LETTER_DIR=/home/ubuntu/judgekit/data/dead-letter
ExecStart=/home/ubuntu/judgekit/judge-worker-rs/target/release/judge-worker
Restart=always
RestartSec=5
# Security hardening
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=true
ReadWritePaths=/home/ubuntu/judgekit/data /home/ubuntu/judgekit/data/dead-letter
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
MemoryMax=4G
CPUQuota=400%

[Install]
WantedBy=multi-user.target
