# scripts/online-judge-backup.service
[Unit]
Description=JudgeKit PostgreSQL backup
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/online-judge
ExecStart=/bin/bash -lc 'mkdir -p /home/ubuntu/backups/online-judge && DB_DIALECT=postgresql /bin/bash /home/ubuntu/online-judge/scripts/backup-db.sh "/home/ubuntu/backups/online-judge/judge-$(date +%%Y%%m%%d-%%H%%M%%S).sql.gz"'
